# habitsync/rootfs/etc/s6-overlay/s6-rc.d/init-habitsync-users/up
/etc/s6-overlay/s6-rc.d/init-habitsync-users/run