# grocy/rootfs/etc/s6-overlay/s6-rc.d/init-keygen/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -e

certfile="/ssl/$(bashio::config 'certfile')"
keyfile="/ssl/$(bashio::config 'keyfile')"

if ! bashio::fs.file_exists "$certfile" || ! bashio::fs.file_exists "$keyfile"; then
    bashio::log.warning "SSL is enabled, but certfile or keyfile is missing. Using self-signed certificate..."

    certfile="/config/keys/cert.crt"
    keyfile="/config/keys/cert.key"

    /usr/local/bin/selfsigned-ssl-gen.sh "$certfile" "$keyfile"
fi
