# owncloud/rootfs/etc/s6-overlay/s6-rc.d/svc-owncloud/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -e

cd /var/www/owncloud || exit 1

exec /usr/bin/entrypoint /usr/bin/owncloud server