# immich_power_tools/rootfs/etc/s6-overlay/s6-rc.d/svc-tools/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

cd /app

# make redis folders
exec node server.js