# .common/postgres/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/command/with-contenv bashio
# shellcheck shell=bash

if ${DELETE_DB:-false} && ${RESTORE_BACKUP:-false}; then
  rm -rf $PGDATA
fi

# Initialize the bundled PostgreSQL data directory on first boot (or skip if
# already initialized). Plain-bash version (no bashio::log.* calls) so this
# works in any base image that ships s6-overlay v3 without a working bashio
# bootstrap. Bash helpers come from /usr/local/bin (postgres-entrypoint.sh)
# via COPY .common/postgres. DB_STORAGE_TYPE tunes HDD/SSD planner.

set -Eeo pipefail

# `.` the helpers under bash rather than via the bashio bootstrap. _is_sourced
# uses FUNCNAME which is bash-only.
# shellcheck source=/dev/null
. /usr/local/bin/postgres-entrypoint.sh

# Set up environment and check if initialization is needed
docker_setup_env
docker_create_db_directories

mkdir -p "/data/logs/postgres"

# permissions
chown postgres:postgres -R "/data/logs/postgres"
chmod 750 -R "/data/logs/postgres"

# Pick a postgresql.conf profile based on DB_STORAGE_TYPE. Default to HDD.
# Uses parameter expansion (`${VAR:-default}`) — which keeps this runnable
# under `set -u` if DB_STORAGE_TYPE is unset.
case "${DB_STORAGE_TYPE:-HDD}" in
  SSD) DB_CFG="postgresql.ssd.conf" ;;
  *)   DB_CFG="postgresql.hdd.conf" ;;
esac
cp -f "/etc/postgresql/${DB_CFG}" /etc/postgresql/postgresql.conf

exec /usr/local/bin/docker-ensure-initdb.sh postgres \
  -c timezone="${TZ:-Etc/UTC}" \
  -c log_timezone="${TZ:-Etc/UTC}" \
  -c config_file=/etc/postgresql/postgresql.conf \
  -c shared_preload_libraries="${POSTGRES_SHARED_PRELOAD_LIBRARIES:-}" \
  -c search_path="${POSTGRES_SEARCH_PATH:-}"
