# immich/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

if bashio::config.true 'delete_db' && bashio::config.true 'restore_backup'; then
  rm -rf $PGDATA
fi

set -Eeo pipefail

# Source the PostgreSQL entrypoint functions
source /usr/local/bin/postgres-entrypoint.sh

# Set up environment and check if initialization is needed
docker_setup_env
docker_create_db_directories

mkdir -p "/data/logs/postgres"

# permissions
chown postgres:postgres -R "/data/logs/postgres"
chmod 750 -R "/data/logs/postgres"

# load postgres config file
case "${DB_STORAGE_TYPE^^}" in
  SSD|HDD)
    bashio::log.info "Using ${DB_STORAGE_TYPE^^} storage"
    cp --update=all "/etc/postgresql/postgresql.${DB_STORAGE_TYPE,,}.conf" "/etc/postgresql/postgresql.conf"
    ;;
  *)
    bashio::log.error "Error: DB_STORAGE_TYPE must be set to 'SSD' or 'HDD'" >&2
    bashio::addon.stop
    exit 1
    ;;
esac

exec /usr/local/bin/docker-ensure-initdb.sh postgres \
    -c timezone="$TZ" \
    -c log_timezone="$TZ" \
    -c config_file="/etc/postgresql/postgresql.conf" &> /dev/null
