# paperless/rootfs/etc/s6-overlay/s6-rc.d/init-db-migration/run
#!/command/with-contenv /usr/bin/bash
set -e

mkdir -p ${PGDATA}-backups
chown postgres:postgres -R ${PGDATA}-backups

if [ -d "${PGDATA}-15-backup" ]; then
    mv $PGDATA "${PGDATA}-backups/$(date +"%Y-%m-%d_%H-%M-%S")"
    mv "${PGDATA}-15-backup" $PGDATA
fi

if [ -f "$PGDATA/PG_VERSION" ] && ! grep -q "$PG_MAJOR" "$PGDATA/PG_VERSION"; then

    backup_root="${PGDATA}-backups"
    target_dir=""

    # Search for all PG_VERSION files and check their content
    while IFS= read -r version_file; do
        if grep -q "$PG_MAJOR" "$version_file"; then
            target_dir=$(dirname "$version_file")
        fi
    done < <(find "$backup_root" -type f -name PG_VERSION 2>/dev/null)

    # Move the last found directory if any were found
    if [ -n "$target_dir" ]; then
        mv "$target_dir" ${PGDATA}
    fi
fi
