# immich/rootfs/etc/s6-overlay/s6-rc.d/db-collation-update/up
/etc/s6-overlay/s6-rc.d/db-collation-update/run