# paperless/rootfs/etc/s6-overlay/s6-rc.d/init-folders/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set +u

declare -r log_prefix="[init-folders]"

declare -r export_dir="${PAPERLESS_EXPORT_DIR:-/usr/src/paperless/export}"
declare -r data_dir="${PAPERLESS_DATA_DIR:-/usr/src/paperless/data}"
declare -r media_root_dir="${PAPERLESS_MEDIA_ROOT:-/usr/src/paperless/media}"
declare -r empty_trash_dir="${PAPERLESS_EMPTY_TRASH_DIR:-/usr/src/paperless/media/.trash}"
declare -r consume_dir="${PAPERLESS_CONSUMPTION_DIR:-/usr/src/paperless/consume}"
declare -r tmp_dir="${PAPERLESS_SCRATCH_DIR:=/tmp/paperless}"

declare -r main_dirs=(
	"${export_dir}"
	"${data_dir}"
	"${media_root_dir}"
	"${empty_trash_dir}"
	"${consume_dir}"
	"${tmp_dir}"
)

declare -r extra_dirs=(
	"${main_dirs[@]}"
	"${data_dir}/index"
	"${media_root_dir}/documents"
	"${media_root_dir}/documents/originals"
	"${media_root_dir}/documents/thumbnails"
)

if [[ -n "${USER_IS_NON_ROOT}" ]]; then
	# Non-root mode: Create directories as current user, warn about permission issues
	echo "${log_prefix} Running in non-root mode, checking directories"
	current_uid=$(id --user)
	current_gid=$(id --group)

	for dir in "${extra_dirs[@]}"; do
		if [[ ! -d "${dir}" ]]; then
			mkdir --parents --verbose "${dir}" || echo "${log_prefix} WARNING: Could not create ${dir} - permission denied"
		fi
		# Check permissions on existing directories too
		if [[ -d "${dir}" && ! -w "${dir}" ]]; then
			echo "${log_prefix} WARNING: No write permission to ${dir}"
		fi
	done

	# Warn about ownership issues
	for dir in "${main_dirs[@]}"; do
		if [[ -d "${dir}" ]]; then
			find "${dir}" -not \( -user ${current_uid} -and -group ${current_gid} \) -exec echo "${log_prefix} WARNING: Permission issue on {}: not owned by current user (${current_uid}:${current_gid})" \; 2>/dev/null || echo "${log_prefix} WARNING: Cannot check permissions on ${dir}"
		fi
	done
else
	# Root mode: Create and fix permissions as needed
	echo "${log_prefix} Running with root privileges, adjusting directories and permissions"

	# First create directories
	for dir in "${extra_dirs[@]}"; do
		if [[ ! -d "${dir}" ]]; then
			mkdir --parents --verbose "${dir}"
		fi
	done

	# Then fix permissions on all directories
	for dir in "${main_dirs[@]}"; do
		find "${dir}" -not \( -user paperless -and -group paperless \) -exec chown --changes paperless:paperless {} +
	done
fi

CONFIG_FILE="/config/paperless.json"

if [ ! -f "$CONFIG_FILE" ]; then
  echo '{}' > "$CONFIG_FILE"
fi

set_config_value() {
  local key="$1"
  local val="$2"
  local tmp
  tmp=$(mktemp) || tmp="${CONFIG_FILE}.tmp"
  jq --arg path "$val" ".${key} = \$path" "$CONFIG_FILE" > "$tmp" && mv "$tmp" "$CONFIG_FILE"
}

handle_folder() {
  local key="$1"
  local new_dir="$2"
  local pretty="$3"

	counter=0

	for dir in "${extra_dirs[@]}"; do
		[[ "${new_dir}" == "${dir}" ]] && coutner=$((counter + 1))
		if [[ "$counter" -ge 2 ]]; then
			bashio::log.error "-----------------------------------------------------------------------------------------------------------------------"
			bashio::log.error "Cannot set two folders to the same path ($new_dir)"
			bashio::log.error "-----------------------------------------------------------------------------------------------------------------------"
			bashio::log.error "Please make sure to set the correct folder paths in the app configuration and restart it"
			bashio::addon.stop
			exit 1
		fi
	done

  # read current from config
  local current
  current=$(jq -r ".${key} // empty" "$CONFIG_FILE")

  # if not set in config, set to provided new_dir
  if [[ -z "$current" ]]; then
    set_config_value "$key" "$new_dir"
    current="$new_dir"
  fi

  # if value changed, handle migration or update
  if [[ "$current" != "$new_dir" ]]; then
    if [[ ! -d "$current" ]]; then
      # previous folder not found — warn and try to accept new_dir if it exists
			bashio::log.warning "-----------------------------------------------------------------------------------------------------------------------"
      bashio::log.warning "Unable to find a ${pretty} in ($current)."
			bashio::log.warning "-----------------------------------------------------------------------------------------------------------------------"
      bashio::log.warning "This is expected if you manually moved the Paperless ${pretty} somewhere else"
      if [[ -d "$new_dir" ]]; then
				bashio::log.info "-----------------------------------------------------------------------------------------------------------------------"
        bashio::log.info "Paperless ${pretty} directory found in ($new_dir)"
				bashio::log.info "-----------------------------------------------------------------------------------------------------------------------"
        bashio::log.info "Setting this as new ${pretty} directory without migration"
      else
				bashio::log.error "-----------------------------------------------------------------------------------------------------------------------"
        bashio::log.error "Unable to find a ${pretty} in ($new_dir)"
				bashio::log.error "-----------------------------------------------------------------------------------------------------------------------"
        bashio::log.error "Paperless cannot start without a valid ${pretty} directory."
        bashio::log.error "Please make sure to set the correct path in the app configuration and restart it"
        bashio::addon.stop
        exit 1
      fi
    else
      # previous folder exists: attempt migration
      /usr/local/bin/storage-folder-migration.sh "$current" "$new_dir"
    fi

    # persist the new path
    set_config_value "$key" "$new_dir"
  fi
}

handle_folder "export_dir"       "${export_dir}"        "export folder"
handle_folder "empty_trash_dir"  "${empty_trash_dir}"   "empty_trash folder"
handle_folder "consume_dir"      "${consume_dir}"       "consume folder"
handle_folder "media_root_dir"   "${media_root_dir}"    "media_root folder"