# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-addon-config/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -e

echo "  Loading env variables:"

OPTIONS_SOURCE=/data/options.json

# Process substitution avoids a subshell so exports persist.
# env_vars array entries are unpacked into individual KEY=value pairs.
while IFS="=" read -r KEY VALUE; do
    [ -z "$KEY" ] && continue

    line="${KEY}=${VALUE}"

    # log redacted config
    case "$KEY" in
        *PASS*|*SECRET*|*KEY*|*TOKEN*)
            echo "    ${KEY}=******"
            ;;
        *)
            echo "    $line"
            ;;
    esac

    # shellcheck disable=SC2163
    export "$line"

    # set .env
    echo "$line" >> /.env || true

    # set /etc/environment
    echo "$line" >> /etc/environment

    # For s6
    if [ -d /var/run/s6/container_environment ]; then
        printf "%s" "$VALUE" > /var/run/s6/container_environment/"${KEY}"
    fi
    echo "export $line" >> ~/.bashrc
done < <(
    jq -r '
        to_entries[]
        | if .key == "env_vars" then
            .value[]? | "\(.key)=\(.value|tostring)"
          else
            select(.value | type != "object" and type != "array")
            | "\(.key)=\(.value|tostring)"
          end
    ' "${OPTIONS_SOURCE}"
)

if [ -n "$TZ" ] && [ -f /etc/localtime ]; then
    if [ -f /usr/share/zoneinfo/"$TZ" ]; then
        echo "Timezone set from $(cat /etc/timezone) to $TZ"
        ln -snf /usr/share/zoneinfo/"$TZ" /etc/localtime
        echo "$TZ" >/etc/timezone
    fi
fi