# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-addon-config/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -e

echo "  Loading env variables:"

OPTIONS_SOURCE=/data/options.json

jq -r 'keys[]' "${OPTIONS_SOURCE}" | while read -r KEY; do
    # export key
    value=$(jq -r --arg key "$KEY" '.[$key]' "${OPTIONS_SOURCE}")

    # Continue for single values
    line="${KEY}=${value}"

    # log redacted config
    case "$KEY" in
        *PASS*|*SECRET*|*KEY*)
            echo "    ${KEY}=******"
            ;;
        *)
            echo "    $line"
            ;;
    esac

    export line

    # set .env
    echo "$line" >> /.env || true

    # set /etc/environment
    echo "$line" >> /etc/environment

    # For s6
    if [ -d /var/run/s6/container_environment ]; then
        printf "%s" "$value" > /var/run/s6/container_environment/"${KEY}"
    fi
    echo "export $line" >> ~/.bashrc
done

if [ -n "$TZ" ] && [ -f /etc/localtime ]; then
    if [ -f /usr/share/zoneinfo/"$TZ" ]; then
        echo "Timezone set from $(cat /etc/timezone) to $TZ"
        ln -snf /usr/share/zoneinfo/"$TZ" /etc/localtime
        echo "$TZ" >/etc/timezone
    fi
fi

ha_address=$(bashio::network.ipv4_address)
ha_address=${ha_address%%/*}

declare protocol="http"
declare certfile="/ssl/$(bashio::config 'certfile')"
declare keyfile="/ssl/$(bashio::config 'keyfile')"

if bashio::config.true 'ssl'; then
    # Require certfile and keyfile
    #bashio::config.require 'certfile'
    #bashio::config.require 'keyfile'

    protocol="https"

    # If certfile or keyfile does not exist, generate self-signed cert
    if ! bashio::fs.file_exists "$certfile" || ! bashio::fs.file_exists "$keyfile"; then
        bashio::log.warning "SSL is enabled, but certfile or keyfile is missing. Using self-signed certificate..."

        certfile="/data/ssl/fullchain.pem"
        keyfile="/data/ssl/privkey.pem"

        /usr/local/bin/selfsigned-ssl-gen.sh "$certfile" "$keyfile"
    fi

    printf "%s" "$keyfile" > /var/run/s6/container_environment/N8N_SSL_KEY
    printf "%s" "$certfile" > /var/run/s6/container_environment/N8N_SSL_CERT

    printf "%s" "$protocol" > /var/run/s6/container_environment/N8N_PROTOCOL
fi

if ! bashio::config.has_value "N8N_HOST"; then
    printf "%s" "$ha_address" > /var/run/s6/container_environment/N8N_HOST
fi

if ! bashio::config.has_value "WEBHOOK_URL"; then
    printf "%s" "${protocol}://${ha_address}:$(bashio::addon.port 5678)" > /var/run/s6/container_environment/WEBHOOK_URL
fi