# grist/rootfs/etc/s6-overlay/s6-rc.d/svc-grist/run
#!/command/with-contenv bashio
# shellcheck shell=bash

cd /grist || exit 1

exec ./sandbox/docker_entrypoint.sh node ./sandbox/supervisor.mjs
