# collabora_code/rootfs/etc/s6-overlay/s6-rc.d/svc-collabora/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

s6-setuidgid cool /start-collabora.sh