# immich_frame/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::net.wait_for "$PORT" localhost 900

bashio::log.info "Starting NGinx..."
exec nginx -c /etc/nginx/nginx.conf -p /etc/nginx
