# immich/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::net.wait_for $SERVER_PORT localhost 900

bashio::log.info "Starting NGinx..."
exec nginx
