# .common/nginx/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -e

bashio::log.info "[svc-nginx] Waiting for service ..."

bashio::net.wait_for "$SERVICE_PORT" localhost 900

bashio::log.info "[svc-nginx] Starting NGinx..."

exec nginx
