# spoolman/rootfs/etc/s6-overlay/s6-rc.d/svc-spoolman/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -e;

chown 1000:1000 -R /home/app/.local/share/spoolman

cd /home/app/spoolman || exit 1

exec /home/app/spoolman/entrypoint.sh