# memos/rootfs/etc/s6-overlay/s6-rc.d/svc-memos/run
#!/command/with-contenv bashio
# shellcheck shell=bash

mkdir -p "$MEMOS_DATA"

bashio::log.info "Starting Memos"

exec /usr/local/memos/memos
