# monica/rootfs/etc/s6-overlay/s6-rc.d/svc-monica/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Start the application
exec entrypoint.sh apache2-foreground