# monica/rootfs/etc/s6-overlay/s6-rc.d/svc-monica/run
#!/command/with-contenv bash
# shellcheck shell=bash

# Start the application
exec entrypoint.sh apache2-foreground