# paperless/rootfs/etc/s6-overlay/s6-rc.d/init-db-vacuum/run
#!/command/with-contenv /usr/bin/bash
set -e

vacuum -u $POSTGRES_USER --all --analyze-in-stages