# memos/rootfs/etc/s6-overlay/s6-rc.d/svc-postgres/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -Eeo pipefail

# Start PostgreSQL server
bashio::log.info "Starting PostgreSQL server..."

s6-notifyoncheck -s 5000 -n 300 -w 1000 -c "pg_isready -U $POSTGRES_USER" \
  postgres-entrypoint.sh postgres \
    -c timezone="$TZ" \
    -c log_timezone="$TZ"
