# jellyfin/rootfs/etc/s6-overlay/s6-rc.d/svc-jellyfin/run
#!/command/with-contenv bashio
# shellcheck shell=bash

cd /jellyfin || exit 1

exec /jellyfin/jellyfin
