# paperless/rootfs/etc/s6-overlay/s6-rc.d/svc-tika/run
#!/command/with-contenv /usr/bin/bash

s6-setuidgid tika exec java -cp "/tika-server-standard-${TIKA_VERSION}.jar:/tika-extras/*" org.apache.tika.server.core.TikaServerCli -h 0.0.0.0 $0 $@