# ezbookkeeping/rootfs/etc/s6-overlay/s6-rc.d/init-addon-config/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -e

echo "  Loading env variables:"

OPTIONS_SOURCE=/data/options.json

jq -r '
    reduce to_entries[] as $item (
        {};
        if $item.value | type == "object" then
        . + $item.value
        else
        . + { ($item.key): $item.value }
        end
    )
    | to_entries[]
    | "\(.key)\u0000\(.value|tostring)"
    ' $OPTIONS_SOURCE | while read -r -d $'\0' key && read -r value; do

    # Continue for single values
    line="$key=$value"

    # log redacted config
    case "$key" in
        *PASS*|*SECRET*|*KEY*)
            echo "    $key=******"
            ;;
        *)
            echo "    $line"
            ;;
    esac

    export line

    # set .env
    echo "$line" >> /.env || true

    # set /etc/environment
    echo "$line" >> /etc/environment

    # For s6
    if [ -d /var/run/s6/container_environment ]; then
        printf "%s" "$value" > /var/run/s6/container_environment/"$key"
    fi
    echo "export $line" >> ~/.bashrc
done

if [ -n "$TZ" ] && [ -f /etc/localtime ]; then
    if [ -f /usr/share/zoneinfo/"$TZ" ]; then
        echo "Timezone set to $TZ"
        ln -snf /usr/share/zoneinfo/"$TZ" /etc/localtime
        echo "$TZ" >/etc/timezone
    fi
fi