# hoarder/rootfs/etc/s6-overlay/s6-rc.d/svc-meilisearch/run
#!/command/with-contenv bashio
# shellcheck shell=bash

mkdir -p $MEILI_DIR

cd $MEILI_DIR

exec /bin/meilisearch